# app/files/alerta/uwsgi.service
[Unit]
Description=uWSGI service for app {{ app }}

[Service]
ExecStart={{ bin_dir }}/uwsgi --ini {{ ini }}

[Install]
WantedBy=multi-user.target
